# src/wsl/docker/supervisord.conf
[supervisord]
nodaemon=true

[program:vnc]
command=/startup.sh
